# kiosk-print-agent/siramatik-print-agent.service
# systemd — tablet acilista print-agent + (gerekirse) git pull
#
# 1) Betik calistirilabilir olsun:
#    chmod +x /home/alper/apps/SIRA_YAZICI/kiosk-print-agent/boot-print-agent.sh
#
# 2) Unit dosyasini kopyala ve etkinlestir:
#    sudo cp siramatik-print-agent.service /etc/systemd/system/
#    sudo systemctl daemon-reload
#    sudo systemctl enable --now siramatik-print-agent.service
#
# Not: User/Group ve ASAGIDAKI yollari kendi kurulumuna gore degistir (monorepo ise
#      .../SIRAMATIK/kiosk-print-agent olabilir).

[Unit]
Description=Siramatik kiosk print-agent (baslangicta kill + git pull + main.py)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=alper
Group=lp
Environment=PRINT_LISTEN_HOST=127.0.0.1
Environment=PRINT_LISTEN_PORT=9131
Environment=PRINT_DEVICE=/dev/usb/lp0
# Ornek: Environment=PRINT_JOB_PREFIX_HEX=1b6501
ExecStart=/bin/sh /home/alper/apps/SIRA_YAZICI/kiosk-print-agent/boot-print-agent.sh
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
